# rootfs/etc/s6-overlay/s6-rc.d/init-download-resources/up
/etc/s6-overlay/s6-rc.d/init-download-resources/run