# rootfs/etc/services.d/lczero-logmon/run
#!/usr/bin/env sh

echo "> LCZero Log Monitor service started"
exec tail -F /lczero/logs/lczero.log
