# rootfs/etc/xinetd.d/lczero
service lczero
{
  disable = no
  server = /lczero/bin/lczero_xinetd
  server_args = --config=/lczero/bin/lc0.config
  socket_type = stream
  protocol = tcp
  user = lczero
  wait = no
  port = 3333
  # only_from = localhost
}
