# rootfs/etc/services.d/lczero/run
#!/usr/bin/env sh

exec 2>&1
exec /bin/s6-tcpserver -u 1000 -g 1000 -v 0.0.0.0 3333 /lczero/lc0_init